# smartinfo19-99.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:
#
#
# This file should be move to /etc/cron.d/smartinfo19-99
# NOTE: no extension!
#
#
# m     h       dom     mon     dow     command
*/5     *       *       *       *       root  ( /home/beheer/ubundiagd/smartinfo19-99.sh |logger -p info -t smartinfo19-99 )
