# temperv14.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:
#
#
# This file should be move to /etc/cron.d/temperv14
# NOTE: no extension!
#
#
# m     h       dom     mon     dow     command
*/3     *       *       *       *       root  ( /home/beheer/ubundiagd/temperv14.sh |logger -p info -t temperv14 )
